# co2mond.service
[Unit]
Description=CO2 Monitor Daemon

[Service]
ExecStart=/usr/local/bin/co2mond -P 0.0.0.0:8080

[Install]
WantedBy=multi-user.target
